# source/etc/rc.d/init.d/ec2-swap
#!/bin/bash
#
# Init file for ec2-swap
#
# chkconfig: 2345 20 20
# description:  ec2-swap is to create swapon swap image on EC2 instance store
#

# source function library
. /etc/rc.d/init.d/functions
. /etc/sysconfig/ec2-swap

RETVAL=0

start() {
  # When system rebooted, return function without any process below.
  if [ -e $SWAPFILENAME ]; then
    logger -i -p local0.notice -t ec2-swap "a swap file exist, exit"
    return 
  fi

	# Get ephemeral disk device name by EC2 meta-data server
	INSTANCEVOL_REALNAME=/dev/`curl -s 169.254.169.254/latest/meta-data/block-device-mapping/$INSTANCEVOL_VIRTNAME`
	# Get memory capacity
	MEMSIZE=`cat /proc/meminfo | grep MemTotal | awk '{print $2}'`

	# Determine swap file size from memory capacity according to below URL
  # https://access.redhat.com/site/documentation/en-US/Red_Hat_Enterprise_Linux/6/html/Installation_Guide/s2-diskpartrecommend-x86.html
	if [ $MEMSIZE -lt 2097152 ]; then
	  SIZE=$((MEMSIZE * 2))k
	elif [ $MEMSIZE -lt 8388608 ]; then
	  SIZE=${MEMSIZE}k
	elif [ $MEMSIZE -lt 67108864 ]; then
	  SIZE=$((MEMSIZE / 2))k
	else
	  SIZE=4194304k
	fi

	# Unmount and convert ephemeral disk from ext3 to ext4 in order to execute fallocate
	if cat /proc/mounts | grep $INSTANCEVOL_MOUNTPOINT; then
		umount $INSTANCEVOL_MOUNTPOINT
	fi
  case "$INSTANCEVOL_CONVMETHOD" in
    mkfs)
      mkfs -t ext4 $INSTANCEVOL_REALNAME
      ;;
    *)
      tune2fs -O extent,flex_bg,huge_file,uninit_bg,dir_nlink,extra_isize $INSTANCEVOL_REALNAME
      e2fsck -fDy $INSTANCEVOL_REALNAME
      ;;
  esac

	# Mount and create swap file
	if [ ! -d $INSTANCEVOL_MOUNTPOINT ]; then
	  mkdir $INSTANCEVOL_MOUNTPOINT
	fi
	mount -t ext4 $INSTANCEVOL_REALNAME $INSTANCEVOL_MOUNTPOINT
  fallocate -l $SIZE $SWAPFILENAME && mkswap $SWAPFILENAME && logger -i -p local0.notice -t ec2-swap "$SWAPFILENAME($SIZE) created"
  umount $INSTANCEVOL_MOUNTPOINT

  # Configure /etc/fstab to enable auto mount ephemeral disk and swapfile
  if ! (cat /etc/fstab | grep $INSTANCEVOL_MOUNTPOINT); then
    echo "$INSTANCEVOL_REALNAME $INSTANCEVOL_MOUNTPOINT auto defaults 0 0" >> /etc/fstab
    logger -i -p local0.notice -t ec2-swap "ephemeral disk entry add in /etc/fstab"
	fi
  if ! (cat /etc/fstab | grep $SWAPFILENAME); then
    echo "$SWAPFILENAME swap swap defaults 0 0" >> /etc/fstab
    logger -i -p local0.notice -t ec2-swap "swap entry add in /etc/fstab"
	fi

  # This script run after rc.sysinit, so it is necessary to mount and swapon explicitly.
  mount -a
  swapon -a
}

stop () {
	:
}

case "$1" in
        start)
                start
                ;;
        stop)
                ;;
        *)
                echo $"Usage: $0 {start}"
                RETVAL=1
esac
exit $RETVAL
